# advanced/active-container/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=tail -f parse.rb
