# advanced/api/supervisord.conf
[supervisord]
nodaemon=true

[program:slapin]
command=rackup -o 0.0.0.0 -p 4700
directory=/api
autorestart=true
