# advanced/api/supervisord.conf
[supervisord]
nodaemon=true

[program:slapin]
command=rackup -p 4700
directory=/api
autorestart=true
