# advanced/api/supervisord.conf
[supervisord]
nodaemon=true

[program:slapin]
command=rackup -p 4567
directory=/api
stdout_logfile=/api/log/%(program_name)s.log
stderr_logfile=/api/log/%(program_name)s.log
autorestart=true
